# systemd/speakip.service
# Speak the wlan0 and eth0 ip addresses at startup
# This requires that espeak is installed and the sound module loaded
# yum install -y espeak
# /sbin/modprobe snd_bcm2835
#
# Place into /usr/lib/systemd/system and run:
# systemctl enable speakip.service
[Unit]
Description=Speak the wlan0 and eth0 ip addresses at startup
After=network-online.target

[Service]
ExecStart=/root/NativeRaspberryPiBeaconParser/bin/run-speakip.sh start
ExecStop=/root/NativeRaspberryPiBeaconParser/bin/run-speakip.sh stop
RemainAfterExit=yes
Type=forking

[Install]
WantedBy=multi-user.target
